# .launch/com.bgem3.rerank.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bgem3.rerank</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>uvicorn</string>
        <string>bgem3_rerank:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8002</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yapilymm/Downloads/projects/bgem3</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/yapilymm/Downloads/projects/bgem3/logs/bgem3_rerank.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yapilymm/Downloads/projects/bgem3/logs/bgem3_rerank.log</string>
</dict>
</plist>